# os_configurations/linux/etc/profile.d/vim.sh
p_=/opt/vim/bin
case ":${PATH:-}:" in
    *:$p_:*) ;;
    *) PATH=$p_:$PATH ;;
esac
export PATH
